# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:lavalink]
command=java -Xmx512m -Xms128m -XX:+UseG1GC -XX:MaxGCPauseMillis=200 -XX:+UnlockExperimentalVMOptions -XX:+UseContainerSupport -Djdk.tls.client.protocols=TLSv1.2 -jar /opt/Lavalink/Lavalink.jar
directory=/opt/Lavalink
user=lavalink
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/lavalink.err.log
stdout_logfile=/var/log/supervisor/lavalink.out.log
priority=1
stopwaitsecs=30

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=2
stopwaitsecs=10
